# init.d/zonegen
#!/bin/sh /etc/rc.common

START=23

start(){
    /bin/zonegen
}
